# src/crontab
# Run dumpsterr every hour
0 * * * * cd /app && /usr/local/bin/python src/main.py
